# udev/90-monitor-plug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/local/bin/toggle_screen.sh"
